# containers/kea/assets/s6-rc.d/kea-ctrl-agent/run
#!/bin/sh
exec 2>&1
exec kea-ctrl-agent -c /etc/kea/kea-ctrl-agent.conf