# container/assets/s6-rc.d/kea-run-prepare/up
if { mkdir -p /run/kea }
if { chown nobody:nogroup /run/kea }
chmod 02755 /run/kea