# containers/kea/assets/s6-rc.d/kea-dhcp6/run
#!/bin/sh
exec 2>&1
exec kea-dhcp6 -c /etc/kea/kea-dhcp6.conf