# container/assets/s6-rc.d/kea-dhcp-ddns/run
#!/bin/sh
exec 2>&1
exec kea-dhcp-ddns -c /etc/kea/kea-dhcp-ddns.conf