# container/assets/s6-rc.d/kea-dhcp4/run
#!/bin/sh
exec 2>&1
exec kea-dhcp4 -c /etc/kea/kea-dhcp4.conf